# viviti-update.timer
[Unit]
Description=Viviti OTA update timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=4h
Unit=viviti-update.service

[Install]
WantedBy=timers.target
